# proxy/macsurf-proxy.service
[Unit]
Description=MacSurf Proxy — TLS-stripping HTTP proxy for Classic Mac OS
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=macsurf
Group=macsurf
EnvironmentFile=/etc/macsurf-proxy/env
ExecStart=/usr/local/bin/macsurf-proxy $MACSURF_PROXY_ARGS
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=macsurf-proxy

# Hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
